# src/teleop_ackermann_joy/config/85-jscal-custom-calibration.rules
ACTION=="add", KERNEL=="js[0-9]*", ATTRS{uniq}=="20423552434B", RUN+="/usr/bin/jscal -s 4,1,0,127,127,4793344,4793344,1,0,127,127,4793344,4793344,1,0,112,142,5534582,5534582,1,0,112,142,5534582,5534582 /dev/input/js%n", SYMLINK+="wfly_et07"
